# code-server/root/etc/s6-overlay/s6-rc.d/init-cabal/run
#!/usr/bin/with-contenv bash

# This script initializes cabal by creating the specified directory,
# optionally setting permissions to 777, and optionally updating the package list.
#
# Environment variables:
#   CABAL_DIR             Directory for storing cabal data (default: /config/.haskell)
#   CABAL_UPDATE_PACKAGES Boolean to control if 'cabal update' should run (default: true)

# Set PUID and PGID from environment or default to 1000
: "${PUID:=1000}"
: "${PGID:=1000}"

# Set default options for cabal behavior if not defined externally
: "${CABAL_DIR:=/config/.haskell}" # Default directory for cabal
: "${CABAL_UPDATE_PACKAGES:=true}" # Enables cabal update by default

# Helper function to determine if a value is "truthy"
is_truthy() {
    local val
    # Return false if variable is not set or empty
    if [ -z "$1" ]; then
        return 1
    fi
    
    # Convert to lowercase
    val=$(echo "$1" | tr '[:upper:]' '[:lower:]')
    
    # Check for truthy values
    case "$val" in
        true|yes|1|on|y)
            return 0 # True
            ;;
        *)
            return 1 # False
            ;;
    esac
}

# Create CABAL_DIR if it doesn't exist and set permissions to 777
if [ ! -d "$CABAL_DIR" ]; then
  echo "Creating directory: $CABAL_DIR"
  mkdir -p "$CABAL_DIR"
fi

# Set ownership of CABAL_DIR to PUID:PGID
echo "Setting ownership of $CABAL_DIR to ${PUID}:${PGID}"
chown -R "${PUID}:${PGID}" "$CABAL_DIR"

if is_truthy "$CABAL_UPDATE_PACKAGES"; then
  if [ ! -d "$CABAL_DIR/packages/hackage.haskell.org" ]; then
    echo "Updating cabal packages"
    s6-setuidgid "${PUID}:${PGID}" cabal update
  else
    echo "Skipping cabal update - package list already exists"
  fi
else
  echo "Skipping cabal update - environment variable is set to false"
fi
