# code-server/root/etc/s6-overlay/s6-rc.d/init-cabal/run
#!/usr/bin/with-contenv bash

CABAL_DIR=~/.haskell

export CABAL_DIR
export PATH="$CABAL_DIR/bin:$PATH"

if [ ! -d "$CABAL_DIR/packages/hackage.haskell.org" ]; then
  cabal update
fi
