# code-server/root/etc/s6-overlay/s6-rc.d/init-cabal/run
#!/usr/bin/with-contenv bash

# Set default options for cabal behavior if not defined externally
: "${CABAL_DIR:=/config/.haskell}" # Default directory for cabal
: "${UPDATE_CABAL_PACKAGES:=true}" # Enables cabal update by default

# Function: is_true
# Description: Returns true if the provided value is a recognized truthy string.
is_true() {
  local value="${1,,}"
  [[ "$value" =~ ^(true|yes|1|on|y)$ ]]
}

# Create CABAL_DIR if it doesn't exist and set permissions to 777
if [ ! -d "$CABAL_DIR" ]; then
  echo "Creating directory: $CABAL_DIR"
  mkdir -p "$CABAL_DIR"
fi

# Set permissions to read, write, and execute for all users
echo "Setting permissions for $CABAL_DIR to read, write, and execute for all users"
chmod -R 777 "$CABAL_DIR"

# Check if cabal update is enabled
if [ ! is_true "$UPDATE_CABAL_PACKAGES" ]; then
  echo "Skipping cabal update"
  exit 0
fi

# Update the package list if it doesn't exist
if [ ! -d "$CABAL_DIR/packages/hackage.haskell.org" ]; then
  echo "Updating cabal packages"
  cabal update
fi
