# code-server/root/etc/s6-overlay/s6-rc.d/svc-tailscale-up/run
#!/usr/bin/with-contenv bash
set -euxo pipefail

# Helper function to determine if a value is "truthy"
is_truthy() {
    local val
    # Return false if variable is not set or empty
    if [ -z "$1" ]; then
        return 1
    fi
    
    # Convert to lowercase
    val=$(echo "$1" | tr '[:upper:]' '[:lower:]')
    
    # Check for truthy values
    case "$val" in
        true|yes|1|on|y)
            return 0 # True
            ;;
        *)
            return 1 # False
            ;;
    esac
}

# Check if Tailscale is enabled
if ! is_truthy "${TAILSCALE_ENABLED}"; then
    echo "INFO: Tailscale is disabled. Exiting."
    exit 0
fi

LOGIN_FLAGS=""

# Provide a Tailscale authentication key if configured
if [ -v TAILSCALE_AUTHKEY ]; then
    echo "INFO: Using Tailscale authkey."
    LOGIN_FLAGS="${LOGIN_FLAGS} --authkey=${TAILSCALE_AUTHKEY}"
fi

# Provide a custom login server if configured
if [ -v TAILSCALE_LOGIN_SERVER ]; then
    echo "INFO: Using custom login server: ${TAILSCALE_LOGIN_SERVER}"
    LOGIN_FLAGS="${LOGIN_FLAGS} --login-server=${TAILSCALE_LOGIN_SERVER}"
fi

if [ -n "${LOGIN_FLAGS}" ]; then
    echo "INFO: Logging in to Tailscale with flags: ${LOGIN_FLAGS}"
    tailscale login ${LOGIN_FLAGS} || true
fi

# Configure a Tailscale hostname if set
# Note: This might affect how other devices on your Tailscale network identify this device
if [ -v TAILSCALE_HOSTNAME ]; then
    echo "INFO: Setting Tailscale hostname to: ${TAILSCALE_HOSTNAME}"
    tailscale set --hostname=${TAILSCALE_HOSTNAME}
fi

# Configure Tailscale to act as an exit node if requested
# Note: This requires approval in the Tailscale admin console
if [ -v TAILSCALE_BE_EXIT_NODE ]; then
    if is_truthy "${TAILSCALE_BE_EXIT_NODE}"; then
        echo "WARNING: Configuring as an exit node. Tailscale admin console approval is required."
        TAILSCALE_BE_EXIT_NODE=true
    else
        TAILSCALE_BE_EXIT_NODE=false
    fi
    tailscale set --advertise-exit-node=${TAILSCALE_BE_EXIT_NODE}
fi

# Accept DNS settings from Tailscale if requested
# Note: Accepting DNS settings may impact internet connectivity
if [ -v TAILSCALE_ACCEPT_DNS ]; then
    if is_truthy "${TAILSCALE_ACCEPT_DNS}"; then
        echo "WARNING: Accepting DNS settings from Tailscale. This may impact internet connectivity."
        TAILSCALE_ACCEPT_DNS=true
    else
        TAILSCALE_ACCEPT_DNS=false
    fi
    tailscale set --accept-dns=${TAILSCALE_ACCEPT_DNS}
fi

# Set the web client preference
# Note: Enabling the web client allows remote management from your Tailnet
if [ -v TAILSCALE_USE_WEBCLIENT ]; then
    if is_truthy "${TAILSCALE_USE_WEBCLIENT}"; then
        echo "WARNING: Enabling the web client for remote management from your Tailnet."
        TAILSCALE_USE_WEBCLIENT=true
    else
        TAILSCALE_USE_WEBCLIENT=false
    fi
    tailscale set --webclient=${TAILSCALE_USE_WEBCLIENT}
fi

# Enable auto-update
# Note: Auto-update ensures Tailscale is always up-to-date
if [ -v TAILSCALE_AUTO_UPDATE ]; then
    if is_truthy "${TAILSCALE_AUTO_UPDATE}"; then
        echo "WARNING: Auto-update is enabled. Tailscale will update automatically."
        TAILSCALE_AUTO_UPDATE=true
    else
        TAILSCALE_AUTO_UPDATE=false
    fi
    tailscale set --auto-update=${TAILSCALE_AUTO_UPDATE}
fi

# Set the Tailscale SSH configuration if provided
# Note: Enabling Tailscale SSH allows SSH access through Tailscale
if [ -v TAILSCALE_USE_SSH ]; then
    if is_truthy "${TAILSCALE_USE_SSH}"; then
        echo "WARNING: Enabling Tailscale SSH. This allows SSH access through Tailscale."
        TAILSCALE_USE_SSH=true
    else
        TAILSCALE_USE_SSH=false
    fi
    tailscale set --ssh=${TAILSCALE_USE_SSH}
fi

tailscale up

# Replace DNS settings with Tailscale DNS if requested
# Note: Replacing DNS settings may impact internet connectivity
if is_truthy "${TAILSCALE_ACCEPT_DNS}" && is_truthy "${TAILSCALE_REPLACE_DNS}"; then
    echo "WARNING: Replacing DNS settings with Tailscale DNS. This may impact internet connectivity."
    python3 /opt/code-server-scripts/modify_resolv_conf.py --nameserver 100.100.100.100
fi

# Configure Tailscale to serve or funnel traffic
# Note: Serving or funneling traffic requires proper configuration
if [ -v TAILSCALE_SERVE_PORT ] && [ -v TAILSCALE_SERVE_MODE ]; then
    if [ -v TAILSCALE_FUNNEL ]; then
        echo "WARNING: Using Tailscale funnel mode, allowing remote Internet access."
        tailscale funnel --"${TAILSCALE_SERVE_MODE}"=443 http://localhost:"${TAILSCALE_SERVE_PORT}" &
    else
        echo "WARNING: Using Tailscale serve mode, allowing remote access from your Tailnet."
        tailscale serve --"${TAILSCALE_SERVE_MODE}"=443 http://localhost:"${TAILSCALE_SERVE_PORT}" &
    fi
fi