# code-server/root/etc/s6-overlay/s6-rc.d/init-patch-json/run
#!/usr/bin/with-contenv bash

# Set default options for patching behavior if not defined externally
: "${PATCH_JSON_USE_FORCE:=true}"           # Enables '--force' flag by default
: "${PATCH_JSON_ENABLE_SORT_KEYS:=true}"    # Enables JSON key sorting by default

# Function: is_true
# Description: Returns true if the provided value is a recognized truthy string.
is_true() {
  local value="${1,,}"
  [[ "$value" =~ ^(true|yes|1|on|y)$ ]]
}

# Function: patch_and_echo
# Description: Applies patches to a JSON file with optional flags, and logs the operation details.
patch_and_echo() {
  local source_file="$1"
  shift
  local patch_files=("$@")
  
  # Build flags based on environment variables
  local flags=""
  if is_true "$PATCH_JSON_USE_FORCE"; then
    flags+=" --force"
  fi
  if ! is_true "$PATCH_JSON_ENABLE_SORT_KEYS"; then
    flags+=" --no-sort-keys"
  fi

  echo "Applying patch to '$source_file' using: ${patch_files[*]}. Flags: '$flags'."
  /opt/shared-venv/bin/python /opt/code-server-scripts/patch_json.py --source "$source_file" $flags --in-place "${patch_files[@]}"
}

# Check if PATCH_JSON is set to a value representing "true"
if [[ -n "$PATCH_JSON" && "${PATCH_JSON,,}" =~ ^(true|yes|1|on|y)$ ]]; then
  patch_and_echo /config/data/Machine/settings.json /usr/share/code-server-scripts/patch_json/*.json
  patch_and_echo /config/data/User/settings.json /usr/share/code-server-scripts/patch_json/*.json
  echo "All JSON patches applied successfully."
else
  echo "PATCH_JSON is not enabled. Skipping JSON patching."
fi
