# code-server/root/etc/s6-overlay/s6-rc.d/fix-venv-ownership/run
#!/usr/bin/with-contenv bash
set -euxo pipefail

# Set PUID and PGID from environment variables or default to 1000
PUID="${PUID:-1000}"
PGID="${PGID:-1000}"

# Set VENV_DIR from environment or default to /opt/shared-venv
VENV_DIR="${VENV_DIR:-/opt/shared-venv}"

# Check if the virtual environment directory exists
if [ ! -d "$VENV_DIR" ]; then
    echo "INFO: Directory $VENV_DIR does not exist. Skipping ownership fix."
    exit 0
fi

echo "INFO: Setting ownership of $VENV_DIR to PUID=$PUID, PGID=$PGID recursively."
if chown -R "$PUID:$PGID" "$VENV_DIR"; then
    echo "INFO: Ownership set successfully."
else
    echo "WARN: Failed to set ownership!"
fi
